# deploy/com.nimbus.server.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.nimbus.server</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/Users/wangqing/sourcecode/agent/agent-framework/nimbus/nimbus</string>
        <string>start</string>
    </array>
    
    <key>WorkingDirectory</key>
    <string>/Users/wangqing/sourcecode/agent/agent-framework/nimbus</string>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>PI_AI_PORT</key>
        <string>3031</string>
        <key>NIMBUS_PORT</key>
        <string>4096</string>
        <key>WEBUI_PORT</key>
        <string>3000</string>
    </dict>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>
    
    <key>StandardOutPath</key>
    <string>/Users/wangqing/sourcecode/agent/agent-framework/nimbus/.logs/launchd-out.log</string>
    
    <key>StandardErrorPath</key>
    <string>/Users/wangqing/sourcecode/agent/agent-framework/nimbus/.logs/launchd-err.log</string>
</dict>
</plist>
